# template-test/test@.path
[Unit]
Description=Path monitoring template for file sync service.
