# template-test/test@.service
[Unit]
Description=Test template for service file with variables

[Service]
ExecStart=echo %i