# template-test/test@.path
[Unit]
Description=Path monitoring file sync service [%t/sync/%i].

[Path]
Unit=test@%i.service
StartLimitIntervalSec=15
MakeDirectory=true
PathModified=%t/sync/%i
