# template-test/test@.service
[Unit]
Description=Test template for service file with variables

[Service]
#EnvironmentFile=/etc/test/test.env
ExecStart=echo %t%i 
