# filesync/filesyncd.service
[Unit]
Description = File syncing between the constellation.
Reqiures = network.service

[Service]
ExecCondition = filesync.sh preflight

