# media-backup/media-backupd.service
[Unit]
Description=Media drive backup daemon
ConditionPathExists=/etc/mediabackupd/mediabackupd.conf

[Service]
Type=simple
CPUQuota=35%
ExecStart=/usr/local/bin/mediabackup
