# media-backup/media-backupd.path
[Unit]
Description=Path monitoring trigger to media-backupd

[Path]
PathModified=/home/orion/Downloads/books
