# media-backup/media-backupd.timer
[Unit]
Description=Timer for 'media-backupd.service'
ConditionPathExists=/mnt/netshare

[Timer]
OnCalendar=Mon *-*-* 00:00:00
